# foot/sv/foot/run
#!/bin/sh -e

exec 2>&1
exec foot --server --log-no-syslog
