# wob/sv/wob/finish
#!/bin/sh -e

. ${XDG_CONFIG_HOME}/wob/conf.sh
rm "${WOBSOCK}"
