# less/profile.d/10-less.sh
if command -v less > /dev/null; then
  export LESS='-FSXiwMRj4a#4z-4'
  export PAGER='less'
fi

