# sv/yambar/run
#!/bin/sh

sv start pipewire

exec 2>&1
exec yambar
